# gateway/ts/script/spool/cron/crontabs/root
#*/1  *  * * *   root    /ts/script/autostart.sh
*/1  *  * * * /ts/script/autostart.sh
